# .config/systemd/user/jacob-hourly.timer
[Unit]
Description=Hourly

[Timer]
OnCalendar=*-*-* *:00:00
Unit=jacob-hourly.service

[Install]
WantedBy=default.target
