# .config/systemd/user/jacob-daily.service
[Unit]
Description=Daily

[Service]
Type=simple
ExecStart=/home/jlewallen/tools/bin/daily.sh

[Install]
WantedBy=default.target
