# .config/systemd/user/jacob-daily.timer
[Unit]
Description=Daily

[Timer]
OnCalendar=*-*-* 00:00:00
Unit=jacob-daily.service

[Install]
WantedBy=default.target
